# contrib/podman/nextcloud-mcp.container
[Unit]
Description=Nextcloud MCP Server
Documentation=https://github.com/pi0n00r/nextcloud-mcp-server

[Container]
ContainerName=nextcloud-mcp
Image=ghcr.io/pi0n00r/nextcloud-mcp-server:v1.1.8
AutoUpdate=registry
PublishPort=127.0.0.1:8000:8000
EnvironmentFile=%h/.config/nextcloud-mcp/env
HealthCmd=curl -fsS http://127.0.0.1:8000/health/live || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthOnFailure=kill
HealthStartPeriod=20s

[Service]
Restart=on-failure
TimeoutStartSec=90

[Install]
WantedBy=default.target
